# debian/resolvconf.upstart
# upstart script for resolvconf

description "Initialize or finalize resolvconf"

start on starting udev

stop on runlevel [06]

pre-start script
	resolvconf --create-runtime-directories
	# Request a postponed update (needed in case the base file has content).
	resolvconf -u
	# Enable updates and perform the postponed update.
	resolvconf --enable-updates
end script

pre-stop script
	resolvconf --disable-updates
end script
