# debian/resolvconf.resolvconf-pull-resolved.path
[Path]
PathChanged=/run/systemd/resolve/resolv.conf
PathExists=/run/systemd/resolve/resolv.conf

[Install]
WantedBy=systemd-resolved.service

